# rfcomm0.path
[Unit]
Description=OBDII port monitoring
# Stops this unit if evnotipi service stopped
BindsTo=evnotipi.service

[Path]
PathExists=/dev/rfcomm0
Unit=evnotipi.service

[Install]
WantedBy=multi-user.target
